# backend/docker/etc/service/roadrunner/run
#!/bin/sh
exec rr serve -w /app/ -c /app/.rr.yaml